# deploy/pqm.service
[Unit]
Description=Pocket Quartermaster (Nuxt SSR)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/pqm
EnvironmentFile=-/opt/pqm/pqm.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/node /opt/pqm/.output/server/index.mjs
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
